# WifiAutologin/Linux/wifi-autologin.service
[Unit]
Description=WifiAutologin daemon
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/WifiAutologin -d
